# dist/src/test/docker-admin/saw-analyze-samples/sample-spark/saw-analyze-spark-samples.service
[Unit]
Description=SAW Analyze Spark samples loader
After=saw-deploy.service

[Service]
Type=oneshot
ExecStart=/root/saw-analyze-samples/sample-spark/load-metric
SyslogIdentifier=saw-analyze-spark-samples

[Install]
WantedBy=saw-deploy.service
