# saw-services/saw-workbench-service/src/main/system/saw-workbench-proxy.socket
[Unit]
Description=SAW Workbench Proxy Socket

[Socket]
ListenStream=9700

[Install]
WantedBy=sockets.target
