# saw-services/saw-workbench-service/src/main/system/sip-workbench-proxy.socket
[Unit]
Description=SIP Workbench Proxy Socket

[Socket]
ListenStream=9700

[Install]
WantedBy=sockets.target
