# sip-services/sip-workbench-executor-service/src/main/system/sip-workbench-executor-proxy.socket
[Unit]
Description=SIP Workbench Executor Proxy Socket

[Socket]
ListenStream=9700

[Install]
WantedBy=sockets.target
