# saw-services/saw-semantic-service/src/main/system/saw-semantic-proxy.socket
[Unit]
Description=SAW Semantic Proxy Socket

[Socket]
ListenStream=127.0.0.1:9900

[Install]
WantedBy=sockets.target
