# dist/src/test/docker/saw-workbench/saw-workbench-testdata.service
[Unit]
Description=SAW Workbench test data loader
After=saw-deploy.service

[Service]
Type=oneshot
ExecStart=/root/saw-workbench/sample-csv/load
SyslogIdentifier=saw-workbench-testdata

[Install]
WantedBy=saw-deploy.service
