# dist/src/test/docker-admin/sip-workbench-samples/sip-workbench-samples.service
[Unit]
Description=SIP Workbench samples loader
After=sip-deploy.service

[Service]
Type=oneshot
ExecStart=/root/saw-workbench-samples/sample-csv/load
SyslogIdentifier=sip-workbench-samples
EnvironmentFile=/etc/bda/secure.env
[Install]
WantedBy=sip-deploy.service
