# saw-services/saw-export-service/src/main/system/saw-export-proxy.socket
[Unit]
Description=SAW Export Proxy Socket

[Socket]
ListenStream=127.0.0.1:9400

[Install]
WantedBy=sockets.target
