# sip-xdf/xdf/bin/templates/xdf-rta.service.template
[Unit]
Description=SIP XDF Real-time Processing Service

[Service]
Type=notify
NotifyAccess=all

ExecStart=/bin/bash -c '<execute-component command with arguments>'

SyslogIdentifier=xdf-rtps
User=mapr
Restart=on-failure
# Consider Java exiting from SIGTERM using exit status 143 to be
# successful
SuccessExitStatus=143