# saw-services/saw-observe-service/src/main/system/saw-observe-proxy.socket
[Unit]
Description=SAW Observe Proxy Socket

[Socket]
ListenStream=9500

[Install]
WantedBy=sockets.target
