# dist/src/test/docker-admin/sip-analyze-samples/sample-spark/sip-analyze-spark-samples.service
[Unit]
Description=SIP Analyze Spark samples loader
After=sip-deploy.service

[Service]
Type=oneshot
ExecStart=/root/saw-analyze-samples/sample-spark/load-metric
SyslogIdentifier=sip-analyze-spark-samples
EnvironmentFile=/etc/bda/secure.env
[Install]
WantedBy=sip-deploy.service
