# dist/src/test/docker/saw-tunnel/saw-tunnel.service
[Unit]
Description=SAW MapR client tunnel
After=network.target

[Service]
# Workaround: The sshuttle documentation provides an example with
# Type=notify, but that seems to cause a time out so leave it out for
# now
#Type=notify
ExecStart=/root/saw-tunnel
Restart=on-failure
SyslogIdentifier=saw-tunnel

[Install]
WantedBy=multi-user.target
