# dist/src/test/docker-admin/saw-analyze-samples/sample-elasticsearch/saw-analyze-elasticsearch-samples.service
[Unit]
Description=SAW Analyze Elasticsearch samples loader
After=saw-deploy.service

[Service]
Type=oneshot
ExecStart=/root/saw-analyze-samples/sample-elasticsearch/load-metric
SyslogIdentifier=saw-analyze-elasticsearch-samples

[Install]
WantedBy=saw-deploy.service
