# dist/src/test/docker-admin/saw-analyze-samples/sample-elasticsearch/sip-analyze-elasticsearch-samples.service
[Unit]
Description=SIP Analyze Elasticsearch samples loader
After=sip-deploy.service

[Service]
Type=oneshot
ExecStart=/root/saw-analyze-samples/sample-elasticsearch/load-metric
SyslogIdentifier=sip-analyze-elasticsearch-samples

[Install]
WantedBy=sip-deploy.service
