# sip-services/sip-metadata-service/src/main/system/sip-metadata-proxy.socket
[Unit]
Description=SIP Metadata Proxy Socket

[Socket]
ListenStream=9500

[Install]
WantedBy=sockets.target
