# saw-services/saw-scheduler-service/src/main/system/saw-scheduler-proxy.socket
[Unit]
Description=SAW Scheduler Proxy Socket

[Socket]
ListenStream=9600

[Install]
WantedBy=sockets.target
