# sip-services/sip-alert-service/src/main/system/sip-alert-proxy.socket
[Unit]
Description=SIP Alert Proxy Socket

[Socket]
ListenStream=9620

[Install]
WantedBy=sockets.target
