# dist/src/test/docker-mapr/sip-mapr-logger@.service
[Unit]
Description=SIP MapR logger: %I
After=default.target

[Service]
ExecStart=/usr/bin/tail -f --retry %I
SyslogIdentifier=sip-mapr-logger@%I
Restart=on-failure

[Install]
WantedBy=default.target
