# saw-services/saw-transport-service/src/main/system/saw-transport-proxy.socket
[Unit]
Description=SAW Transport Proxy Socket

[Socket]
ListenStream=9200

[Install]
WantedBy=sockets.target
