# saw-services/saw-batch-ingestion-service/src/main/system/saw-batch-ingestion-proxy.socket
[Unit]
Description=SAW Batch Ingestion Socket

[Socket]
ListenStream=9952

[Install]
WantedBy=sockets.target
