# saw-services/saw-batch-ingestion-service/src/main/system/saw-batch-ingestion-proxy.socket
[Unit]
Description=SAW Batch Ingestion Socket

[Socket]
ListenStream=127.0.0.1:9903

[Install]
WantedBy=sockets.target
