# dist/src/test/docker/saw-metrics/sample-elasticsearch/saw-analyze-elasticsearch-testdata.service
[Unit]
Description=SAW Analyze Elasticsearch test data loader
After=saw-deploy.service

[Service]
Type=oneshot
ExecStart=/root/saw-metrics/sample-elasticsearch/load-metric
SyslogIdentifier=saw-analyze-elasticsearch-testdata

[Install]
WantedBy=saw-deploy.service
